# darwin/mac-bridge/com.jbeda.mac-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jbeda.mac-bridge</string>

  <!-- inetd-style: launchd listens on the socket, accepts each connection,
       and runs the program once per connection with the socket wired to
       stdin/stdout. No long-running daemon, no polling. -->
  <key>inetdCompatibility</key>
  <dict>
    <key>Wait</key>
    <false/>
  </dict>

  <key>Sockets</key>
  <dict>
    <key>Listeners</key>
    <dict>
      <!-- Bound to loopback: only the SSH reverse tunnel (which lands on the
           Mac's localhost) can reach it. -->
      <key>SockNodeName</key>
      <string>127.0.0.1</string>
      <key>SockServiceName</key>
      <string>17603</string>
      <key>SockType</key>
      <string>stream</string>
      <key>SockFamily</key>
      <string>IPv4</string>
    </dict>
  </dict>

  <!-- Read one tab-separated request line ("verb TAB args") from the socket and
       dispatch it. Only two verbs are honored, each with strict validation, so a
       stray local connection can't run arbitrary commands.
         open URL      : only http/https, then run `open URL`
         code HOST PATH: HOST must be a bare ssh-alias and PATH absolute, then
                         run VS Code's `code` CLI with `remote ssh-remote+HOST PATH` -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>tab=$(printf '\t'); IFS=$tab read -r verb a1 a2 || exit 0; case "$verb" in open) case "$a1" in http://*|https://*) exec /usr/bin/open "$a1";; esac ;; code) case "$a1" in ''|*[!A-Za-z0-9._-]*) exit 0;; esac; case "$a2" in /*) ;; *) exit 0;; esac; PATH="/usr/local/bin:/opt/homebrew/bin:$HOME/.local/bin:$PATH"; exec code --remote "ssh-remote+$a1" "$a2" ;; esac</string>
  </array>
</dict>
</plist>
